# systemd/home-network-duplicate-scan.service
[Unit]
Description=Home Network weekly duplicate report scan
Documentation=file:/home/jellyfish/repo/home-network/docs/operations/storage-monitoring.md
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/mnt/2TB /mnt/4TB /opt/docker

[Service]
Type=oneshot
ExecStart=/opt/docker/appdata/storage-monitoring/run-duplicate-scan.sh
ExecStartPost=/opt/docker/appdata/storage-monitoring/export-prometheus-metrics.sh
Nice=15
IOSchedulingClass=idle
