# systemd/mqtt-backup-discord-bridge.service
[Unit]
Description=Home Network MQTT backup event Discord bridge
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Environment=MQTT_HOST=jellyhome
Environment=MQTT_PORT=1883
Environment=MQTT_TOPIC=home-network/backups/+/borgmatic/event
Environment=MQTT_USERNAME=borgmatic
Environment=MQTT_PASSWORD_FILE=/etc/home-network/mqtt_borgmatic_password
Environment=HERMES_BIN=/home/jellybot/.hermes/hermes-agent/venv/bin/hermes
Environment=SEND_START=0
Environment=SEND_SUCCESS=0
EnvironmentFile=-/etc/home-network/mqtt-backup-discord-bridge.env
ExecStartPre=/usr/bin/test -r /etc/home-network/mqtt_borgmatic_password
ExecStart=/opt/docker/bin/mqtt-backup-discord-bridge
Restart=always
RestartSec=10s
User=jellybot
Group=jellybot
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/jellybot/.hermes

[Install]
WantedBy=multi-user.target
